# profile.d/custom.sh
alias l='ls -lah'
alias vi='/usr/bin/vim'
